# commands/scripts/crontab/hermes-crontab
# Hermes maintenance and quantitative decision-loop schedule

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Minute guard: holdings + persisted recommendations/anchor ETFs. The runner is
# fail-closed unless confirmed positions, MiniQMT quotes, and trade-calendar
# confirmation are all available. Live orders still require daily authorization.
30-59 9 * * 1-5 cd /home/ly/.hermes/research-assistant && .venv_quant/bin/python3 commands/scripts/decision_guard_once.py >> ~/.hermes/logs/decision_guard.log 2>&1
* 10 * * 1-5 cd /home/ly/.hermes/research-assistant && .venv_quant/bin/python3 commands/scripts/decision_guard_once.py >> ~/.hermes/logs/decision_guard.log 2>&1
0-30 11 * * 1-5 cd /home/ly/.hermes/research-assistant && .venv_quant/bin/python3 commands/scripts/decision_guard_once.py >> ~/.hermes/logs/decision_guard.log 2>&1
* 13-14 * * 1-5 cd /home/ly/.hermes/research-assistant && .venv_quant/bin/python3 commands/scripts/decision_guard_once.py >> ~/.hermes/logs/decision_guard.log 2>&1
0 15 * * 1-5 cd /home/ly/.hermes/research-assistant && .venv_quant/bin/python3 commands/scripts/decision_guard_once.py >> ~/.hermes/logs/decision_guard.log 2>&1
55 8 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/vnext_passlist.py >> ~/.hermes/logs/decision_guard.log 2>&1
20 15 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/decision_postmarket.py >> ~/.hermes/logs/decision_review.log 2>&1
30 15 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/decision_certify.py >> ~/.hermes/logs/decision_certification.log 2>&1
20 16 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/generate_capability_status.py >> ~/.hermes/logs/capability_status.log 2>&1
30 17 * * 1-5 cd /home/ly/.hermes/research-assistant && scripts/backup_data_to_d.sh >> ~/.hermes/logs/data_backup.log 2>&1
*/30 8-18 * * 1-5 cd /home/ly/.hermes/research-assistant && scripts/data_recovery_guard.sh >> ~/.hermes/logs/data_recovery_guard.log 2>&1
10 16 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/vnext_event_truth_fetch.py >> ~/.hermes/logs/vnext_event_truth.log 2>&1
40 15 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/vnext_realized_collect.py >> ~/.hermes/logs/vnext_realized.log 2>&1

# L2 is a digest; L3/L4 are delivered immediately by the guard.
35 11 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands /usr/bin/python3 -m factor_lab.decision_loop.cli flush-l2 >> ~/.hermes/logs/decision_guard.log 2>&1
5 15 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands /usr/bin/python3 -m factor_lab.decision_loop.cli flush-l2 >> ~/.hermes/logs/decision_guard.log 2>&1
