# commands/scripts/crontab/hermes-crontab
# Hermes maintenance and quantitative decision-loop schedule

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Minute guard: holdings + persisted recommendations/anchor ETFs. The runner is
# fail-closed unless confirmed positions, MiniQMT quotes, and trade-calendar
# confirmation are all available. Live orders still require daily authorization.
# Canonical live market truth has one ingestion owner. The non-blocking global
# lock prevents overlap with backup, recovery, and post-market DataHub writers.
* 9-10 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/datahub_live_snapshot.py >> ~/.hermes/logs/datahub_live_snapshot.log 2>&1
0-30 11 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/datahub_live_snapshot.py >> ~/.hermes/logs/datahub_live_snapshot.log 2>&1
* 13-14 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/datahub_live_snapshot.py >> ~/.hermes/logs/datahub_live_snapshot.log 2>&1
0 15 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/datahub_live_snapshot.py >> ~/.hermes/logs/datahub_live_snapshot.log 2>&1
30-59 9 * * 1-5 cd /home/ly/.hermes/research-assistant && .venv_quant/bin/python3 commands/scripts/decision_guard_once.py >> ~/.hermes/logs/decision_guard.log 2>&1
* 10 * * 1-5 cd /home/ly/.hermes/research-assistant && .venv_quant/bin/python3 commands/scripts/decision_guard_once.py >> ~/.hermes/logs/decision_guard.log 2>&1
0-30 11 * * 1-5 cd /home/ly/.hermes/research-assistant && .venv_quant/bin/python3 commands/scripts/decision_guard_once.py >> ~/.hermes/logs/decision_guard.log 2>&1
* 13-14 * * 1-5 cd /home/ly/.hermes/research-assistant && .venv_quant/bin/python3 commands/scripts/decision_guard_once.py >> ~/.hermes/logs/decision_guard.log 2>&1
0 15 * * 1-5 cd /home/ly/.hermes/research-assistant && .venv_quant/bin/python3 commands/scripts/decision_guard_once.py >> ~/.hermes/logs/decision_guard.log 2>&1
* 8-20 * * * cd /home/ly/.hermes/research-assistant && source ~/.config/hermes/runtime.env && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/decision_notification_worker.py >> ~/.hermes/logs/decision_notifications.log 2>&1
55 8 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/run_scheduled_dag.py premarket >> ~/.hermes/logs/premarket_dag.log 2>&1
# Post-market write DAG uses dependencies and checkpoints instead of assuming a
# prior job completed because its wall-clock slot elapsed. Independent status
# and backup jobs still run after a research failure.
45 15 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/run_scheduled_dag.py postmarket >> ~/.hermes/logs/postmarket_dag.log 2>&1
*/30 8-18 * * 1-5 cd /home/ly/.hermes/research-assistant && scripts/data_recovery_guard.sh >> ~/.hermes/logs/data_recovery_guard.log 2>&1
0 10 * * 0 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 commands/scripts/run_scheduled_dag.py weekly_datahub >> ~/.hermes/logs/datahub_cron.log 2>&1

# L2 is a digest; L3/L4 use the durable outbox worker above.
35 11 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 -m factor_lab.decision_loop.cli flush-l2 >> ~/.hermes/logs/decision_guard.log 2>&1
5 15 * * 1-5 cd /home/ly/.hermes/research-assistant && PYTHONPATH=commands .venv_quant/bin/python3 -m factor_lab.decision_loop.cli flush-l2 >> ~/.hermes/logs/decision_guard.log 2>&1
